# systemd/nvml-user-exporter.service
[Unit]
Description=NVML per-user GPU usage Prometheus exporter
After=network-online.target nvidia-persistenced.service
Wants=network-online.target

[Service]
Type=simple
# Assumes the repo has been copied to /opt/nvml_exporter and a venv has been
# created at /opt/nvml_exporter/.venv with deps installed (see README "Systemd").
# If you keep the repo under /home/<user>/, change both paths below and remove
# ProtectHome=true.
ExecStart=/opt/nvml_exporter/.venv/bin/python /opt/nvml_exporter/nvml_user_exporter.py --port 9835 --addr 0.0.0.0
Restart=on-failure
RestartSec=5

# Running as root is the safe default: resolving /proc/<pid>/status for
# *other* users' PIDs may be blocked on distros that set hidepid=2 on /proc,
# in which case 'nobody' would tag everything as <gone>. If your /proc is
# permissive (hidepid=0), you can switch to User=nobody Group=nogroup for
# defense in depth.
User=root

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
